# assets/systemd/aistack-idle.timer
[Unit]
Description=aistack Idle Evaluation Timer
Documentation=https://github.com/polygonschmiede/aistack
Requires=aistack-idle.service

[Timer]
# Run every 10 seconds
OnBootSec=30s
OnUnitActiveSec=10s

# Persistent across reboots
Persistent=true

[Install]
WantedBy=timers.target
