# assets/udev/70-aistack-wol.rules
# Persist Wake-on-LAN configuration via aistack CLI
# Triggers on physical network interfaces and reapplies configured WoL mode.
SUBSYSTEM=="net", ACTION=="add|change", ATTR{type}=="1", KERNEL!="lo", RUN+="/usr/local/bin/aistack wol-apply %k"
